# deploy/systemd/goodos-backup-retention.timer
[Unit]
Description=Run GoodOS backup retention cleanup daily

[Timer]
OnCalendar=*-*-* 02:45:00
Persistent=true
Unit=goodos-backup-retention.service

[Install]
WantedBy=timers.target
