# deploy/systemd/goodbase-telemetry-retention.timer
[Unit]
Description=Run Goodbase telemetry retention nightly

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=900
Unit=goodbase-telemetry-retention.service

[Install]
WantedBy=timers.target
