# deploy/systemd/goodbase-observability.service
[Unit]
Description=Goodbase production observability controller
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/var/www/Goodbase/deploy/observability
Environment=OBSERVABILITY_ENV_FILE=/etc/goodbase/observability.env
ExecStartPre=/var/www/Goodbase/deploy/observability/scripts/validate.sh
ExecStart=/usr/bin/docker compose --env-file /etc/goodbase/observability.env up -d --remove-orphans
ExecReload=/usr/bin/docker compose --env-file /etc/goodbase/observability.env up -d --remove-orphans
ExecStop=/usr/bin/docker compose --env-file /etc/goodbase/observability.env stop
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
