# deploy/systemd/goodbase-data-platform.service
[Unit]
Description=Goodbase REST, GraphQL, pooling, and Realtime data platform
Requires=docker.service postgresql.service
After=docker.service postgresql.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/var/www/Goodbase/deploy/data-platform
EnvironmentFile=/etc/goodbase/data-platform.env
ExecStart=/usr/bin/docker compose --env-file /etc/goodbase/data-platform.env up -d --remove-orphans postgrest pgbouncer-transaction pgbouncer-session realtime
ExecReload=/usr/bin/docker compose --env-file /etc/goodbase/data-platform.env up -d --remove-orphans postgrest pgbouncer-transaction pgbouncer-session realtime
ExecStop=/usr/bin/docker compose --env-file /etc/goodbase/data-platform.env stop postgrest pgbouncer-transaction pgbouncer-session realtime
TimeoutStartSec=180
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
