# deploy/systemd/goodbase-backup-retention.service
[Unit]
Description=Goodbase encrypted backup retention
After=postgresql.service goodos-enterprise-backup.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/sbin/goodbase-enterprise-retention
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/backups/goodos-enterprise-v1
